# files/.config/systemd/user/rescuetime.service
[Unit]
Description=rescuetime

[Service]
ExecStart=/usr/bin/rescuetime

[Install]
WantedBy=default.target
