# files/.config/systemd/user/wiki.service
[Unit]
Description=personal wiki

[Service]
ExecStart=/home/daniel/bin/wiki

[Install]
WantedBy=default.target
